# .crontab
#
# $Id: .crontab 1894 2007-03-30 20:27:03Z KevinAnLipscomb $
#
MAILTO=kevin.lipscomb@kvrs.org
30 13 * * * wget -q http://www.oscarsystem.com/daily-avail-report.phtml
55 01 23 * * wget -q http://www.oscarsystem.com/trends.phtml
30 06 * * 1 wget -q http://www.emsi.org/lamp-app/lms-applicant/weekly-applications-report.php
