# trunk/oscar/.crontab
#
# $Id$
#
MAILTO=kevin.lipscomb@kvrs.org
30 13 * * * wget http://www.oscarsystem.com/daily-avail-report.phtml
55 01 23 * * wget http://www.oscarsystem.com/trends.phtml
30 06 * * 1 wget http://www.emsi.org/lamp-app/lms-applicant/weekly-applications-report.phtml
