# trunk/oscar/.crontab
#
# $Id$
#
MAILTO=kevin.lipscomb@kvrs.org
#
20 03 * * * lynx http://www.oscarsystem.com/daily-avail-report.phtml
